# crontab
30 21 * * 1,2,3,4 python3 create_issue.py >> /var/log/cron_create.log 2>&1
30 22 * * 1,2,3,4 python3 close_issue.py >> /var/log/cron_close.log 2>&1

