# crontab
* * * * * python3 create_issue.py >> /var/log/cron_create.log 2>&1
* * * * * python3 close_issue.py >> /var/log/cron_close.log 2>&1
