# crontab
* * * * * python3 create_issue.py >> /var/log/cron_create.log 2>&1
